# platform/desktop/linux/packages/core/base/systemd/units/foo.automount
[Unit]
Description=Awesome Unit

# Automount at boot
[Automount]
Where=
DirectoryMode=

[Install]
WantedBy=multi-user.target
